# conf.d/english.conf
[program:english]

command=/usr/local/bin/ffmpeg_english.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
startsecs=5
startretries=600
