# conf.d/espanol.conf
[program:espanol]

command=/usr/local/bin/ffmpeg_espanol.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
startsecs=5
startretries=600
